# deploy/rs3-discord-bot.service
[Unit]
Description=RS3-inspired Discord Bot
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=tj
Group=tj
WorkingDirectory=/home/tj/rs3-bot
Environment=NODE_ENV=production
EnvironmentFile=/home/tj/rs3-bot/.env
ExecStart=/home/tj/.nvm/versions/node/v22.23.1/bin/node /home/tj/rs3-bot/dist/index.js
Restart=always
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=full

[Install]
WantedBy=multi-user.target
